# crontab
*/5 * * * * bash /home/ubuntu/enzyme-explorer/check-upstream.sh
