# builder/crontab
# m h  dom mon dow   command
*/5 * * * * cronic bash /app/check-upstream.sh
