# builder/crontab
# m h  dom mon dow   command
*/5 * * * * root /bin/bash /app/check-upstream.sh
