# builder/crontab
# m h  dom mon dow   command
*/5 * * * * /app/check-upstream.sh
